# images/klutshnik-sbox-rpi/input/setup_init.sh
#!/sbin/openrc-run
# shellcheck shell=ash

depend() {
  need localmount
  need networking
}

start() {
	klutshnik-setup
}
